# cli/samweb.target
[Unit]
Description=Samantha's http and chat servers
Wants=samhttp.service samchat.service

[Install]
WantedBy=graphical.target
